# www-apps/redmine/files/redmine.socket
[Unit]
Description=Redmine listen Socket

[Socket]
ListenStream=/var/run/redmine_fcgi.socket
Accept=false

[Install]
WantedBy=sockets.target
